# scripts/node-security-notifier/com.hidari.node-security-notifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hidari.node-security-notifier</string>
    <key>ProgramArguments</key>
    <array>
        <string>__DOTFILES_DIR__/scripts/node-security-notifier/run.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/.local/state/node-security-notifier/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.local/state/node-security-notifier/launchd.log</string>
</dict>
</plist>
